# deploy/solar-display.service
# Solar Display systemd service
# Install: sudo cp solar-display.service /etc/systemd/system/
# Enable: sudo systemctl enable solar-display
# Start: sudo systemctl start solar-display

[Unit]
Description=Solar Display Kiosk System
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
User=kz
WorkingDirectory=/opt/solar-display
EnvironmentFile=-/opt/solar-display/.env
Environment=NODE_ENV=production
Environment=DATA_DIR=/opt/solar-display/data
Environment=LOG_DIR=/opt/solar-display/logs
Environment=PORT=3000
SyslogIdentifier=solar-display
ExecStart=/usr/bin/node apps/server/dist/server.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/solar-display/data /opt/solar-display/logs /opt/solar-display/uploads/images /opt/solar-display/uploads/brand

[Install]
WantedBy=multi-user.target
