# deploy/solar-display.service
# Solar Display systemd service
# Install: sudo cp solar-display.service /etc/systemd/system/
# Enable: sudo systemctl enable solar-display
# Start: sudo systemctl start solar-display

[Unit]
Description=Solar Display Kiosk System
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
User=pi
WorkingDirectory=/data/solar-display
EnvironmentFile=-/data/solar-display/.env
Environment=NODE_ENV=production
Environment=DATA_DIR=/data/solar-display/data
Environment=LOG_DIR=/data/solar-display/logs
Environment=PORT=3000
Environment=KIOSK_USER=pi
SyslogIdentifier=solar-display
ExecStart=/usr/bin/node apps/server/dist/server.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/data/solar-display/data /data/solar-display/logs /data/solar-display/uploads/images /data/solar-display/uploads/brand

[Install]
WantedBy=multi-user.target
